# openhost/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:sculptor-backend]
# Sculptor backend (packaged PyInstaller binary) on the external port (8080).
# Runs as devuser via the user directive.
command=/opt/sculptor/sculptor_backend/sculptor_backend --no-open-browser
directory=/opt/sculptor
user=%(ENV_SCULPTOR_RUN_USER)s
environment=HOME="%(ENV_OPENHOST_APP_DATA_DIR)s/home"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:seed-workspace]
# One-shot: create a starter git repo and register it as a Sculptor workspace.
command=python3 /app/setup/seed_workspace.py
directory=/app
user=%(ENV_SCULPTOR_RUN_USER)s
environment=HOME="%(ENV_OPENHOST_APP_DATA_DIR)s/home"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
